# src/fxkeyd.service
[Unit]
Description=key remapping daemon (fxkeyd)

[Service]
Type=simple
ExecStart=/usr/local/bin/fxkeyd

[Install]
WantedBy=multi-user.target
